# mpvpaper-live-wallpaper/cpu-power-limits.service
[Unit]
Description=Apply balanced CPU power limits on boot
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /home/yourusername/.scripts/game-mode.sh off

[Install]
WantedBy=multi-user.target
